# docs/SUBMIT_2026-07-21/06_runtime/twoder-status.service
[Unit]
Description=2DER status regeneration (one shot)

[Service]
Type=oneshot
User=takasan
ExecStart=/srv/workcell/status/run-status.sh
